# oardocker/templates/common/images/frontend/systemd/oardocker-apache.socket
[Unit]
Description=HTTPD Socket

[Socket]
ListenStream=80

[Install]
WantedBy=oardocker.target
