# oardocker/templates/common/images/frontend/systemd/oardocker-nginx.socket
[Unit]
Description=Nginx reverse proxy Socket

[Socket]
ListenStream=80

[Install]
WantedBy=oardocker.target
