# oardocker/templates/common/images/frontend/systemd/oardocker-php.socket
[Socket]
ListenStream=/var/run/php5-fpm.sock
SocketUser=www-data
SocketGroup=www-data

[Install]
WantedBy=sockets.target
