# images/server/supervisor/oar-server.conf
[program:oar-server]
command=/sbin/%(program_name)s
process_name=%(program_name)s
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
