# kameleon/steps/data/lib/systemd/system/oar-docker.target
#  This file is part of oar-docker project.

[Unit]
Description=Basic oardocker System
