# oardocker/templates/images/base/supervisor/rsyslogd.conf
[program:rsyslogd]
command=/usr/sbin/rsyslogd -n
process_name=%(program_name)s
autostart=true
autorestart=true
user=root
directory=/
priority=912
redirect_stderr=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
