# oardocker/templates/images/base/supervisor/sshd.conf
[program:sshd]
command=/usr/local/sbin/sshd
process_name=%(program_name)s
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
