# oardocker/templates/images/server/supervisor/oar-server.conf
[program:oar-server]
command=/usr/local/sbin/oar_server_cmd
process_name=%(program_name)s
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
