# oardocker/templates/images/server/supervisor/postgresql.conf
[program:postgresql]
command=/usr/local/sbin/setuser postgres /usr/lib/postgresql/9.1/bin/postmaster -D /var/lib/postgresql/9.1/main/ -c config_file=/etc/postgresql/9.1/main/postgresql.conf
process_name=%(program_name)s
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
priority=10
