# images/services/supervisor/dnsmasq.conf
[program:dnsmasq]
command=/usr/local/sbin/dnsmasq_cmd
process_name=%(program_name)s
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
